# systemd/dropbox-report.service
[Unit]
Description=Drop Box 12-hour summary report

[Service]
Type=oneshot
ExecStart=/opt/dropbox/report.sh
